# supervisor/create-topics.ini
[program:topic]
command=/usr/bin/create-kafka-topics.sh
autostart=true
autorestart=false