# supervisor/initialize.ini
[program:initialize]
command=/usr/bin/start-all.sh restart
autostart=true
autorestart=false