# supervisor/kafka.ini
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=false
autorestart=false