# supervisor/schemaregistry.ini
[program:schemaregistry]
command=/usr/bin/start-schemaregistry.sh
autostart=false
autorestart=false