# supervisor/zookeeper.ini
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
autostart=false
autorestart=false
stdout_logfile = /var/log/zookeeper/zookeeper.log
stdout_logfile_backups = 0
redirect_stderr = true
